# config/systemd/vod-api.service
[Unit]
Description=VOD Processor API
After=network.target

[Service]
Type=simple
User=vodprocessor
Group=vodprocessor
WorkingDirectory=/opt/vod-processor
Environment="PATH=/opt/vod-processor/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/vod-processor/venv/bin/python /opt/vod-processor/api_server.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/vod-processor/api.log
StandardError=append:/var/log/vod-processor/api.error.log

[Install]
WantedBy=multi-user.target
